# data/etc/hotplug.d/online/50-opkg-restore
if [ ! -f /etc/opkg-restore ]; then exit 0; fi
if [ -f /sbin/extendroot ] && [ -f /etc/extendroot ] && [ -f /etc/extendroot.failed ]; then exit 0; fi
if [ -f /sbin/extendroot ] && [ -f /etc/extendroot ] && [ ! -f /etc/extendroot.done ]; then
	if /sbin/extendroot; then touch /etc/extendroot.done; reboot; exit 0; else touch /etc/extendroot.failed; exit 0; fi
fi

OPKG_CONF="init main"
for OPKG_CONF in ${OPKG_CONF}
do
	if [ ! -e /etc/opkg-restore-"${OPKG_CONF}" ] && lock -n /var/lock/opkg-restore && opkg update; then
		. /etc/profile.d/opkg.sh
		opkg restore "${OPKG_CONF/main}" 2>&1 | logger -t opkg
		touch /etc/opkg-restore-"${OPKG_CONF}"
		lock -u /var/lock/opkg-restore
		reboot
		break
	fi
done
